# factory/propolis/smf/serial.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
    <service name='site/buildomat/serial' type='service' version='0'>
        <create_default_instance enabled='true'/>
        <single_instance/>

        <dependency name='multi-user' grouping='require_all'
          restart_on='none' type='service'>
            <service_fmri value='svc:/milestone/multi-user'/>
        </dependency>

        <exec_method name='start' type='method' timeout_seconds='60'
          exec='/bin/bash /vm/serial.sh'/>
        <exec_method name='stop' type='method' timeout_seconds='60'
          exec=':kill'/>

        <property_group name='startd' type='framework'>
            <propval name='duration' type='astring' value='child'/>
        </property_group>

        <stability value='Unstable'/>

        <template>
            <common_name>
                <loctext xml:lang='C'>serial data forwarding</loctext>
            </common_name>
            <description>
                <loctext xml:lang='C'>serial data forwarding</loctext>
            </description>
        </template>

    </service>
</service_bundle>
